# node/supervisord.conf
[supervisord]
nodaemon=true

[program:elasticsearch]
command=/startNode.sh
autorestart=true
priority=30

